# packaging/gsus.service
[Unit]
Description=gsus daemon (user)

[Service]
Type=simple
ExecStart=%h/.local/bin/gsusd
Restart=on-failure

[Install]
WantedBy=default.target
